# scraper.conf
start on runlevel [2345]
stop on runlevel [016]

respawn
setuid ubuntu
setgid ubuntu
chdir /home/ubuntu/cl-rental
env WORKDIR=/home/ubuntu/cl-rental
exec python3 scrape_loop.py