# deployment/supervisord.conf
[supervisord]
nodaemon=true

[program:afinder]
directory=/opt/wwc/apartment-finder
command=python3 main_loop.py
