# deployment/supervisord.conf
[supervisord]
nodaemon=true

[program:afinder]
user=ubuntu
directory=/opt/wwc/cl-rental
command=python3 main_loop.py
